# download_agent/momentum-sync.timer
[Unit]
Description=Momentum wallpaper - changes 4x daily (morning, noon, evening, night)

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 20:00:00
OnCalendar=*-*-* 23:00:00
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
