# download_agent/momentum-sync.timer
[Unit]
Description=Momentum wallpaper - changes 3x daily (morning, afternoon, night)

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 13:00:00
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
